# launchd/com.skylar.artifactd-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.skylar.artifactd-tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/cloudflared</string>
    <string>tunnel</string>
    <string>--config</string>
    <string>/Users/skylarpayne/.cloudflared/artifacts.yml</string>
    <string>run</string>
    <string>d8df3cd1-438e-47c9-996f-2c5d2b271bf1</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/skylarpayne/Library/Logs/artifactd/cloudflared.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/skylarpayne/Library/Logs/artifactd/cloudflared.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/skylarpayne</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
