# launchd/com.skylar.artifactd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.skylar.artifactd</string>

  <key>UserName</key>
  <string>skylarpayne</string>

  <key>GroupName</key>
  <string>staff</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/skylarpayne/artifactd &amp;&amp; export ARTIFACTD_COOKIE_SECRET=&quot;$(cat /Users/skylarpayne/.hermes/artifacts/.cookie-secret)&quot; &amp;&amp; exec /Users/skylarpayne/artifactd/.venv/bin/artifactd --home /Users/skylarpayne/.hermes/artifacts --public-base-url https://artifacts.skylarbpayne.com serve --host 127.0.0.1 --port 8787</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/skylarpayne/artifactd</string>

  <key>StandardOutPath</key>
  <string>/Users/skylarpayne/Library/Logs/artifactd/artifactd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/skylarpayne/Library/Logs/artifactd/artifactd.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/skylarpayne</string>
    <key>PATH</key>
    <string>/Users/skylarpayne/artifactd/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>ARTIFACTD_HOME</key>
    <string>/Users/skylarpayne/.hermes/artifacts</string>
    <key>ARTIFACTD_PUBLIC_BASE_URL</key>
    <string>https://artifacts.skylarbpayne.com</string>
  </dict>
</dict>
</plist>
